# deploy/supervisor.conf
# /etc/supervisor/conf.d/watchboard.conf
#
# After copying: supervisorctl reread && supervisorctl update

# ── Laravel Horizon (queue worker) ────────────────────────────────────────────
[program:watchboard-horizon]
process_name=%(program_name)s
command=php /var/www/watchboard/artisan horizon
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/watchboard/storage/logs/horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
stopsignal=SIGTERM
stopasgroup=true
killasgroup=true

# ── Laravel Reverb (WebSocket server on 127.0.0.1:8080) ───────────────────────
[program:watchboard-reverb]
process_name=%(program_name)s
command=php /var/www/watchboard/artisan reverb:start --host=0.0.0.0 --port=8080 --no-interaction
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/watchboard/storage/logs/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopsignal=SIGTERM
stopasgroup=true
killasgroup=true

# ── Group (start/stop both together) ─────────────────────────────────────────
[group:watchboard]
programs=watchboard-horizon,watchboard-reverb
